# misc/systemd/net-porter.service
[Unit]
Description=net-porter service
Documentation=https://github.com/a-light-win/net-porter

[Service]
Type=simple
ExecStart=/usr/lib/netavark/net-porter server -c /etc/net-porter/config.json
Restart=on-failure
RestartSec=10s

# Security hardening
ProtectSystem=strict
ReadWritePaths=/run /tmp
# NoNewPrivileges=true blocks setns() which CNI plugins require
NoNewPrivileges=false
ProtectKernelTunables=true
ProtectControlGroups=true
# Allow-list: only net and mnt namespace operations are permitted.
# Worker needs: setns(CLONE_NEWNS), unshare(CLONE_NEWNS), and
# CNI plugins need setns(CLONE_NEWNET).
# Using allow-list (no ~ prefix) for explicit control.
RestrictNamespaces=net mnt

[Install]
WantedBy=multi-user.target
