# misc/systemd/net-porter.service
[Unit]
Description=net-porter service
Documentation=https://github.com/a-light-win/net-porter

[Service]
Type=simple
ExecStart=/usr/lib/netavark/net-porter server -c /etc/net-porter/config.json
Restart=on-failure
RestartSec=10s

# Security hardening
ProtectSystem=strict
ReadWritePaths=/run
# NoNewPrivileges=true blocks setns() which CNI plugins require
# RestrictNamespaces=~net limits namespace operations to net only
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=~net

[Install]
WantedBy=multi-user.target
